# dbus-saver/polonium-saver.service
[Unit]
Description=DBus saver for Polonium

[Service]
Type=dbus
ExecStart=polonium-saver
BusName=xyz.vaughanm.polonium
